# configs/sds-controller.service
[Unit]
Description=SDS Controller
Documentation=https://github.com/liliang-cn/sds
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sds
ExecStart=/opt/sds/bin/sds-controller --config /etc/sds/controller.toml
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/sds

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sds-controller

[Install]
WantedBy=multi-user.target
